# backend/openpos-ki.service
[Unit]
Description=OpenPositioning KI-Steuerung
After=network.target openpos.service ollama.service
Wants=openpos.service ollama.service

[Service]
Type=simple
User=moritz
WorkingDirectory=/home/moritz/openpos/backend
EnvironmentFile=/home/moritz/openpos/backend/.env
ExecStart=/usr/bin/python3 /home/moritz/openpos/backend/ki_controller.py
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
